# build/centos-base/container-files/etc/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
user=root
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
minfds=65535
minprocs=65535

[program:crond]
command=/usr/sbin/crond -n
user=root
numprocs=1
startsecs=7
startretries=3
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[include]
files = supervisord.d/*.ini
