# build/kafka/container-files/etc/supervisord.d/kafka.ini
[program:kafka]
command=/opt/kafka_2.12-2.0.0/bin/kafka-server-start.sh  /opt/kafka_2.12-2.0.0/config/server.properties
stdout_logfile=/opt/kafka_2.12-2.0.0/logs/kafka_out.log
stderr_logfile=/opt/kafka_2.12-2.0.0/logs/kafka_err.log
user=root
startsecs=8
startretries=3
autostart=true
autorestart=false
