# build/logstash/container-files/etc/supervisord.d/logstash.ini
[program:logstash]
command=/opt/logstash/bin/logstash -f /opt/logstash/conf/file.conf
user=root
startsecs=8
startretries=3
autostart=true
autorestart=false
