# build/kibana/container-files/etc/supervisord.d/kibana.ini
[program:kibana]
command=sh /opt/kibana-6.6.1-linux-x86_64/bin/kibana
user=root
startsecs=8
startretries=3
autostart=true
autorestart=false